# certmonger-exporter.service.in
[Service]
Type=notify
User=nobody
ExecStart=/usr/bin/python3 -I @libexecdir@/certmonger-exporter.pyz
Restart=always

[Install]
WantedBy=multi-user.target
