# certmonger-exporter.service.in
[Unit]
Documentation=man:certmonger-exporter(8)

[Service]
Type=notify
User=nobody
ExecStart=${python} -I ${libexecdir}/certmonger-exporter.pyz
Restart=always

[Install]
WantedBy=multi-user.target
